# bitsd-scripts/etc/crontabs/root

MAILTO=""

0  8 * * * /usr/bin/bitsd-start.sh
0 23 * * * /usr/bin/bitsd-stop.sh
